# qubes-mgmt-salt-user/templates/provides-onlykey-agent/files/vm/etc/systemd/user/gpg-agent.service.d/10-ok-gpg-agent.conf
[Service]
ExecStart=/usr/bin/ok-gpg-agent

# vim: set ft=systemd :
