# qubes-mgmt-salt-user/common/hardening/files/usr/lib/systemd/system/hardened_malloc.service
[Unit]
Description="Configure hardened_malloc linker preload"

[Service]
ExecStart=/usr/share/qubes-user/preload-hardened-malloc

[Install]
WantedBy=multi-user.target

