# src/salt/templates/provides-flatpak/files/systemd/user-flatpak-apps.path
# vim: set ft=systemd :
[Unit]
Description=Watcher for changes to Flatpak apps

[Path]
PathModified=%h/.local/share/flatpak/exports/applications
Unit=qubes-sync-appmenu.service
TriggerLimitIntervalSec=1h

[Install]
WantedBy=default.target

