# src/salt/templates/provides-flatpak/files/systemd/user-update-flatpaks.timer
# vim: set ft=systemd :
[Unit]
Description=Periodically check for Flatpak updates


[Timer]
OnBootSec=5min
OnUnitActiveSec=2d

[Install]
WantedBy=timers.target
