# qubes-mgmt-salt-user/templates/provides-socks5-ivpn/files/vm/usr/lib/systemd/system/ivpn-tun2proxy.service.j2
[Unit]
Description=Check connection status against an IVPN endpoint
After=network.target

[Service]
Type=oneshot
ExecStart={{ svc_bin_path }} start

[Install]
WantedBy=multi-user.target

{#- vim: set syn=systemd.jinja : -#}
