# srv/user/salt/templates/provides-flatpak/files/update-flatpak-user.service
[Unit]
Description=Update user flatpaks
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak --user update -y
User=user

